# tools/rootfs/alpine/patch/etc/init.d/diskresize
#!/sbin/openrc-run

description="disk resize"

start_pre() {
	if [ ! -e /var/lib/misc/firstrun ]; then
		ebegin "disk resize begin"
		/usr/bin/resize-helper
		mkdir -p /var/lib/misc
		touch /var/lib/misc/firstrun
		eend 0
	fi
}