# tools/rootfs/alpine/patch/etc/init.d/linkmount
#!/sbin/openrc-run

description="linkmount"

start_pre() {
	/etc/init.d/S20linkmount start
}